# php-7.0-fpm-alpine/rootfs/etc/supervisor.d/supervisord.ini
[supervisord]
nodaemon = true
pidfile = /var/run/supervisord/supervisord.pid
logfile = /dev/null

[program:sshd]
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
command = /usr/sbin/sshd -D -e

[program:php-fpm]
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
command= /usr/local/sbin/php-fpm -O
